# wabrum_bot.service
[Unit]
Description=Wabrum Content Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wabrum_bot
ExecStart=/home/ubuntu/wabrum_bot/venv/bin/python main.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/wabrum_bot/.env

[Install]
WantedBy=multi-user.target
